# config/picoclaw.service
[Unit]
Description=PicoClaw AI Agent
After=network-online.target piclaw-first-boot.service
Wants=network-online.target
ConditionPathExists=/opt/piclaw/.setup-complete

[Service]
Type=simple
User=picoclaw
Group=picoclaw
ExecStart=/opt/picoclaw/picoclaw gateway
WorkingDirectory=/home/picoclaw/.picoclaw/workspace
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# PicoClaw reads config from ~/.picoclaw/config.json
Environment=HOME=/home/picoclaw

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/home/picoclaw
ReadWritePaths=/opt/picoclaw
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
